# meta-rcar-gen3/recipes-bsp/optee/optee-client/optee.service
[Unit]
Description=optee services

[Service]
ExecStart=/usr/bin/tee-supplicant
ExecStop=/usr/bin/killall -s KILL tee-supplicant

[Install]
WantedBy=multi-user.target
